# cal-notify.service
[Unit]
Description = Calendar notifications

[Service]
Type = oneshot
WorkingDirectory = %h/cal-notify/
ExecStart = /bin/sh cal-notify.sh
